# com.antigravity.titans.maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.titans.maintenance</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/yoyocubano/Desktop/École - trabajo /PIF_Documentation/.agent/titans_system/maintenance.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/yoyocubano/Desktop/École - trabajo /PIF_Documentation/.agent/titans_system/logs/maintenance.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yoyocubano/Desktop/École - trabajo /PIF_Documentation/.agent/titans_system/logs/maintenance_err.log</string>
</dict>
</plist>
